# system/supervisor/frontend.conf
[program:test]
directory="project directory"
command=npm start
autorestart=true 
stderr_logfile=/var/log/test.err.log
stdout_logfile=/var/log/test.out.log