# system/supervisor/frontend.conf
[program:frontend]
directory="project directory"
command=npm start
autostart = true
autorestart=true 
stderr_logfile=/var/log/test.err.log
stdout_logfile=/var/log/test.out.log